# debian/postrm
#!/bin/bash -e

if [[ "$1" != "upgrade" ]]; then
  echo "Disabling FlashMQ systemd service"
  systemctl disable flashmq.service || echo "Ignoring..."
fi
